# files/tmp/crontab
0 0 * * * /bin/gc_builders >> /data/builder-gc/builder_gc.log 2>&1

